# email_classifier_brain/email-classifier.service
[Unit]
Description=Email Classification Service
After=network.target

[Service]
User=email_bot
Group=email_bot
WorkingDirectory=/home/email_bot/email_classifier_brain
ExecStart=/home/email_bot/email_classifier_brain/run_service.sh
Restart=always
RestartSec=5
EnvironmentFile=/home/email_bot/email_classifier_brain/.env

[Install]
WantedBy=multi-user.target
